# ops/test/supervisord.scheduler_ingest_test.conf
; Test-only supervisord for scheduler_ingest in docker-compose.test.yml.
;
; Runs ONLY the ingestor program. The prod supervisord also starts router,
; feature_extractor, and stats_aggregator; we skip all three because:
;
;   - router would consume from /data/ipc/crawl_result/crawler_NN/, which
;     this test does not populate (the patrol worker writes directly to
;     ingestor_NN/ following the IPC contract).
;   - feature_extractor processes content-bearing records; sitemap-discovered
;     URLs are status="new" with no content.
;   - stats_aggregator writes to summary_daily / domain_stats_daily; not
;     needed for the sitemap test and would just produce unused rows.

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[program:ingestor]
directory=/app
process_name=%(program_name)s_%(process_num)02d
command=python -m containers.scheduler_ingest.ingestor.main --config /app/containers/scheduler_ingest/config/ingest.test.yaml --ingestor-id %(process_num)d
numprocs=16
numprocs_start=0
; OMP / BLAS pin matches the prod supervisord — the ingestor module loads
; numpy/scipy at import time. We disable inline scoring so the artifact
; isn't required, but keep the BLAS pin to avoid thread thrash.
environment=PYTHONPATH="/app",OMP_NUM_THREADS="1",OPENBLAS_NUM_THREADS="1",MKL_NUM_THREADS="1"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
