# containers/scheduler_control/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[program:offerer]
directory=/app
process_name=%(program_name)s_%(process_num)02d
command=python -m containers.scheduler_control.offerer.main --config /app/containers/scheduler_control/config/control.yaml --offerer-id %(process_num)d
numprocs=16
numprocs_start=0
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:golden_discovery_ranker_v1]
directory=/app
process_name=%(program_name)s_%(process_num)02d
command=python -m containers.scheduler_control.scorer.main --config /app/containers/scheduler_control/config/control.yaml --worker-id %(process_num)d
numprocs=16
numprocs_start=0
environment=PYTHONPATH="/app"
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:accounting_rolloff]
directory=/app
command=python -m containers.scheduler_control.accounting.main --config /app/containers/scheduler_control/config/control.yaml
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
